# dev-db/riak/files/init.d/riak
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need logger
}

start() {
	ebegin "Starting riak"
	start-stop-daemon --start --name riak --chdir /var/lib/riak --user riak --exec /usr/sbin/riak -- start
	eend $?
}

stop() {
	ebegin "Stopping riak"
	start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 --user riak --exec ${ERTS_PATH}/beam.smp
	rval=$?

	#cleanup: also kill epmd daemon
	sleep 5
	${ERTS_PATH}/epmd -kill &> /dev/null

	eend $rval
}

reload() {
	ebegin "Reloading riak"
	start-stop-daemon --start --name riak --chdir /var/lib/riak --user riak --exec /usr/sbin/riak -- restart
	eend $?
}
